# repo_fetch.cron
# Run repo fetch every 10 minutes
*/10 * * * * root /usr/src/redmine/bin/rails runner "Repository.fetch_changesets" -e production >> /var/log/cron.log 2>&1
#
